# wifi_portal/wifi-portal.service
[Unit]
Description=Raspberry Pi WiFi Configuration Portal
After=network.target wpa_supplicant.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wifi-portal
Environment="PATH=/opt/wifi-portal/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/wifi-portal"
Environment="WIFI_INTERFACE=wlan0"
Environment="WPA_SUPPLICANT_CONF=/etc/wpa_supplicant/wpa_supplicant.conf"
ExecStart=/opt/wifi-portal/venv/bin/python -m uvicorn wifi_portal.webui.app:app --host 0.0.0.0 --port 8080 --log-level info
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
